# bash/network/ssh_fw.service.example
# see details from https://www.raspberrypi.org/documentation/linux/usage/systemd.md

[Unit]
Description=Remote SSH Login Service
After=systemd-networkd-wait-online.service
Requires=systemd-networkd-wait-online.service

[Service]
Type=oneshot
ExecStart=/home/pi/bin/my-awesome-scripts/bash/network/ssh_fw.sh ljishen@piha.soe.ucsc.edu
RemainAfterExit=true
StandardOutput=inherit
StandardError=inherit
User=pi

[Install]
WantedBy=multi-user.target
